# ops/systemd/radar-schedules.timer
[Unit]
Description=Dispara o radar-schedules de hora em hora

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
